# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ============================================================
; Scheduler ARGOS — collecte + trading toutes les 15 min
; ============================================================
[program:scheduler]
command=python -m argos.main run-paper
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
redirect_stderr=false
stdout_logfile=/app/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/scheduler.err.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/app/src",TRADING_MODE="paper",ENABLE_REAL_TRADING="false",DATABASE_PATH="/data/argos.db"

; ============================================================
; Dashboard Streamlit — accessible sur le port 8501
; ============================================================
[program:dashboard]
command=python -m streamlit run src/argos/dashboard/app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
    --server.enableXsrfProtection=false
    --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=false
stdout_logfile=/app/logs/dashboard.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/dashboard.err.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/app/src",TRADING_MODE="paper",ENABLE_REAL_TRADING="false",DATABASE_PATH="/data/argos.db"

; ============================================================
; Watchdog — init DB au démarrage et vérifications initiales
; ============================================================
[program:init]
command=python -m argos.main init-db
directory=/app
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/app/logs/init.log
stderr_logfile=/app/logs/init.err.log
environment=PYTHONPATH="/app/src",DATABASE_PATH="/data/argos.db"
